# crontab
* * * * * root __SCRIPT_DIR__ capture 2>&1 | /usr/bin/logger -t TLB -i
