# crontab
* * * * * cjaure /home/cjaure/timelapsebot/tlb capture 2>&1
