# crontab
* * * * * root /home/cjaure/timelapsebot/tlb capture 2>&1 | /usr/bin/logger -t TLB -i
