# timelapsebot.service
[Unit]
Description=Timelapse bot photo taking

[Service]
Type=oneshot
User=chip
ExecStart=/usr/local/bin/tlb capture

